# desfollow-limpeza.service
[Unit]
Description=Desfollow Limpeza Automática de Jobs
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/desfollow
ExecStart=/usr/bin/python3 /root/desfollow/sistema_limpeza_automatica.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Configurações de segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/desfollow

[Install]
WantedBy=multi-user.target 